# myapp.conf
[program:myapp]

command  = /usr/local/bin/gunicorn myapp:myapp

directory = /root/myapp

user = root

autostart = true

stdout_logfile = /tmp/myapp.log

stderr_logfile = /tmp/myapp_error.log